# crontab
0 1    * * *   root    azure_blob_backup