# crontab
0 1 * * * root /usr/local/bin/azure_blob_backup